# llama-server.service
# llama-server.service — Persistent LLM inference backend for Christopher
#
# Install on t3610:
#   sudo cp llama-server.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable llama-server
#   sudo systemctl start llama-server
#
# Check status:
#   sudo systemctl status llama-server
#   journalctl -u llama-server -f
#
# Christopher then connects with --no-server flag:
#   python3 ~/Projects/Christopher-AI/christopher.py --chat --no-server

[Unit]
Description=llama-server — Local LLM inference for Christopher AI
After=network.target
Documentation=https://gitlab.com/JRM-FusionAL/Christopher-AI

[Service]
Type=simple
User=jrm_fusional
Group=jrm_fusional

ExecStart=/home/jrm_fusional/llama.cpp/build/bin/llama-server \
    -m /home/jrm_fusional/llama.cpp/models/Llama-3.2-3B-Instruct-Q4_K_M.gguf \
    -ngl 99 \
    -t 4 \
    -c 2048 \
    --host 127.0.0.1 \
    --port 8080 \
    --log-disable

# Restart automatically if it crashes
Restart=on-failure
RestartSec=10

# Give the GPU time to be ready on boot
TimeoutStartSec=120

# Log to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llama-server

[Install]
WantedBy=multi-user.target

